# K3-LEDE/files/etc/init.d/k3-txpower
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /bin/sh -c '
		sleep 20
		for iface in wlan0 wlan1; do
			iwconfig "$iface" txpower 20 2>/dev/null || true
		done
	'
	procd_close_instance
}
