# packages/luci-app-aw1k-led/root/etc/hotplug.d/net/31-aw1000-phone-led
#!/bin/sh

[ "$ACTION" = "add" ] || [ "$ACTION" = "remove" ] || exit 0

/usr/bin/aw1000-phone-led-sync.sh
