# packages/rg500q-at-webserver/files/etc/uci-defaults/90-rg500q-at-webserver
#!/bin/sh

chmod 755 /etc/init.d/rg500q-at-webserver 2>/dev/null
/etc/init.d/rg500q-at-webserver enable

exit 0
