# packages/luci-app-modemwebui/root/etc/uci-defaults/90-modemwebui
#!/bin/sh

if [ -x /etc/init.d/modemwebui ]; then
	/etc/init.d/modemwebui enable
	/etc/init.d/modemwebui restart
fi

exit 0
